# backend/com.acai.backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.acai.backend</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/anaconda3/bin/python3.11</string>
        <string>-m</string>
        <string>uvicorn</string>
        <string>src.main:app</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>8000</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/ximenruixue/Desktop/AC_AI/backend</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/ximenruixue/Desktop/AC_AI/backend</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>